# data/mobile-power-saver.service.in
[Unit]
Description=Mobile Power Saver

[Service]
Type=exec
ExecStart=@BIN_DIR@/mobile-power-saver

[Install]
WantedBy=gnome-session-initialized.target
